# supervisord.conf
[supervisord]

[program:start-magpie]
directory = /root
command = /root/deploy-magpie.sh
redirect_stderr = true
stdout_logfile = /root/magpie.log
